# docker/dockerizing/app.conf
[program:app]
command=/usr/bin/dotnet /var/www/app/ExpInvoiceVsc.dll
directory=/var/www/app
autostart=true
autorestart=true
environment=ASPNETCORE_ENVIRONMENT=Production
user=www-data
stopsignal=INT